# openhab.conf
[program:openhab]
command=/openhab/start.sh
process_name=%(program_name)s
user=openhab
