# socat-zwave.conf
[program:socat]
command=/usr/bin/socat -d -d pty,link=/dev/ttyNET0,raw,user=openhab,group=dialout,mode=777 tcp:192.168.10.13:3333
user=openhab
autorestart=true
